# system-daemon/customer-dashboard-updater.service
[Unit]
Description=Customer Dashboard Auto-Update Daemon
After=network-online.target docker.service

[Service]
WorkingDirectory=/opt/customer-dashboard/system-daemon
ExecStart=/usr/bin/node daemon.js
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
