# system-daemon/customer-dashboard-updater.service
[Unit]
Description=Customer Dashboard Auto Updater
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/customer-dashboard/system-daemon
ExecStart=/usr/bin/node /opt/customer-dashboard/system-daemon/daemon.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
